# deploy/supervisor/tv-polymarket-bot.conf
[program:tv-polymarket-bot]
directory=/opt/tradingview-polymarket-bot
command=/opt/tradingview-polymarket-bot/.venv/bin/python /opt/tradingview-polymarket-bot/main.py
autostart=true
autorestart=true
startsecs=5
stderr_logfile=/opt/tradingview-polymarket-bot/logs/supervisor.err.log
stdout_logfile=/opt/tradingview-polymarket-bot/logs/supervisor.out.log
environment=PYTHONUNBUFFERED="1"
